# roles/desktop/files/.config/systemd/user/battery.service
[Unit]
Description=Battery alert

[Service]
Type=oneshot
ExecStart=%h/bin/dunst/battery.sh

[Install]
WantedBy=xorg.target
