# roles/desktop/files/.config/systemd/user/xcompmgr.service
[Unit]
Description=Simple compositing manager

[Service]
ExecStart=/usr/bin/xcompmgr -F

[Install]
WantedBy=xorg.target
