# roles/desktop/files/.config/systemd/user/xmonad.service
[Unit]
Description=Xmonad window manager

[Service]
ExecStart=/usr/bin/xmonad

[Install]
WantedBy=xorg.target
