# roles/docker/files/docker_prune.timer
[Unit]
Description=Docker prune timer

[Timer]
OnCalendar=daily
Unit=docker_prune.service

[Install]
WantedBy=timers.target
