# roles/disk_extra/files/backup_mirror.timer
[Unit]
Description=Backup mirror timer

[Timer]
OnCalendar=daily
Unit=backup_mirror.service

[Install]
WantedBy=timers.target
