# roles/keepassxc/files/.config/autostart/keepassxc.desktop
[Desktop Entry]
Name=KeepassXC
Exec=/usr/bin/keepassxc
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=tru
