# roles/influxdb/files/backup_influxdb.timer
[Unit]
Description=Backup Influxdb timer

[Timer]
OnCalendar=monthly
Unit=backup_influxdb.service

[Install]
WantedBy=timers.target
