# roles/desktop/files/.config/autostart/terminal.desktop
[Desktop Entry]
Name=KeepassXC
Exec=xfce4-terminal --maximize
Terminal=false
Type=Application
StartupNotify=false
