# roles/desktop/templates/80-backlight.rules.j2
RUN+="/bin/chgrp video /sys/class/backlight/{{ host.desktop.acpi }}/brightness"
RUN+="/bin/chmod g+w /sys/class/backlight/{{ host.desktop.acpi }}/brightness"
